# vagrant/upstart/rqworker.conf
start on runlevel [2345]
stop on runlevel [!2345]

respawn
chdir /home/vagrant
exec bin/ralph rqworker ext --settings=ralph.settings.dev
